# config/debian-non-US/cron.hourly
#! /bin/sh
#
# Executed hourly via cron, out of troup's crontab.

set -e
export SCRIPTVARS=/org/non-us.debian.org/dak/config/vars
. $SCRIPTVARS

cd $masterdir
dak import-users-from-passwd
