# config/debian/cron.hourly
#! /bin/sh
#
# Executed hourly via cron, out of troup's crontab.

set -e
set -u
export SCRIPTVARS=/srv/ftp.debian.org/dak/config/debian/vars
. $SCRIPTVARS

date -u > $ftpdir/project/trace/ftp-master.debian.org
dak import-users-from-passwd
dak queue-report -n > $webdir/new.html
dak show-new > /dev/null
