# config/debian/cron.hourly
#! /bin/sh
#
# Executed hourly via cron, out of troup's crontab.

set -e
set -u
export SCRIPTVARS=/org/ftp.debian.org/dak/config/debian/vars
. $SCRIPTVARS

dak import-users-from-passwd
dak queue-report -n > $webdir/new.html
