# config/debian-security/cron.hourly
#! /bin/bash
#
# Executed hourly via cron, out of dak's crontab.

exec /srv/security-master.debian.org/dak/config/debian-security/cronscript hourly
